# contrib/ccpeek.service
[Unit]
Description=ccpeek - Claude Code Chat History Viewer
Documentation=https://github.com/adsorce/ccpeek
After=network.target

[Service]
Type=simple
ExecStart=%h/ccpeek/ccpeek --no-browser --port 8888
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
